# contrib/mapeced.service
[Unit]
Description=MAP-E CE (Customer Edge) daemon
Documentation=https://github.com/negineri/mapeced
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/sbin/mapeced
Restart=on-failure
RestartSec=5s

# /run/mapeced/ を起動時に作成
RuntimeDirectory=mapeced
RuntimeDirectoryMode=0755

# /var/lib/mapeced/ を作成（DUID ファイル保存先）
StateDirectory=mapeced
StateDirectoryMode=0700

[Install]
WantedBy=multi-user.target
